# com.paypalnewsreader.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.paypalnewsreader</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/Users/weili/Desktop/Claude Code Test/paypal-news-reader/server/index.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/weili/Desktop/Claude Code Test/paypal-news-reader</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/paypal-news-reader.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/paypal-news-reader-error.log</string>
</dict>
</plist>
